# packaging/omnios/zoneweaver-agent-smf.xml
<?xml version='1.0'?>
<!DOCTYPE service_bundle SYSTEM '/usr/share/lib/xml/dtd/service_bundle.dtd.1'>
<!--
    SMF manifest for Zoneweaver Agent
-->
<service_bundle type='manifest' name='zoneweaver-agent'>
  <service name='system/virtualization/zoneweaver-agent' type='service' version='1'>
    <create_default_instance enabled='false'/>
    <single_instance/>
    
    <!-- Dependencies -->
    <dependency name='filesystems' grouping='require_all' restart_on='none' type='service'>
      <service_fmri value='svc:/system/filesystem/local:default'/>
    </dependency>
    
    <dependency name='network' grouping='require_all' restart_on='none' type='service'>
      <service_fmri value='svc:/milestone/network:default'/>
    </dependency>
    
    <!-- Start method -->
    <exec_method name='start' type='method' exec='/opt/zoneweaver-agent/startup.sh' timeout_seconds='60'>
      <method_context working_directory='/opt/zoneweaver-agent'>
        <method_credential user='zwagent' group='zwagent' 
            privileges='basic,net_privaddr,sys_admin,proc_owner,file_dac_read,file_dac_write,proc_priocntl' />
        <method_environment>
          <envvar name='PATH' value='/opt/ooce/bin:/opt/ooce/node-22/bin:/usr/gnu/bin:/usr/bin:/usr/sbin:/sbin'/>
          <envvar name='NODE_ENV' value='production'/>
          <envvar name='CONFIG_PATH' value='/etc/zoneweaver-agent/config.yaml'/>
          <envvar name='HOME' value='/var/lib/zoneweaver-agent'/>
        </method_environment>
      </method_context>
    </exec_method>
    
    <!-- Stop method: the :kill token sends SIGTERM to every process in the
         service contract (node plus any zadm vnc / zlogin / host-shell PTY
         children), so the contract empties in ~1s instead of svc.startd having
         to wait out this timeout and force-kill surviving children. Node's own
         SIGTERM handler (lib/Shutdown.js) closes the database cleanly on exit. -->
    <exec_method name='stop' type='method' exec=':kill' timeout_seconds='30' />
    
    <!-- Refresh method -->
    <exec_method name='refresh' type='method' exec=':kill -HUP' timeout_seconds='30' />
    
    <!-- Stability -->
    <stability value='Evolving' />
    
    <!-- Service description -->
    <template>
      <common_name>
        <loctext xml:lang='C'>Zoneweaver Agent</loctext>
      </common_name>
      <description>
        <loctext xml:lang='C'>
          Node.js API backend for Zoneweaver Agent zone hypervisor management.
          Provides RESTful API for managing Bhyve virtual machines, networking,
          storage, and system monitoring on OmniOS/illumos.
        </loctext>
      </description>
      <documentation>
        <manpage title='zoneweaver-agent' section='8' />
        <doc_link name='Zoneweaver Agent Project'
                  uri='https://github.com/Makr91/zoneweaver-agent' />
      </documentation>
    </template>
  </service>
</service_bundle>
